# kamailio-rtpengine/rootfs/etc/services.d/rtpengine/run
#!/usr/bin/with-contenv bash
set -euo pipefail

read_env() {
  [ -f "/run/s6/container_environment/$1" ] && cat "/run/s6/container_environment/$1" || true
}

IP="$(read_env ADVERTISED_IP)"
[ -n "${IP:-}" ] && [ "$IP" != "0.0.0.0" ] || { echo "FATAL: ADVERTISED_IP non valido"; exit 1; }

exec rtpengine \
  --interface="${IP}" \
  --listen-ng=127.0.0.1:22222 \
  --port-min=49000 --port-max=50000 \
  --log-level=4 \
  --timeout=60 \
  --dtls-passive \
  --table=-1 \
  --nftables-chain=
