# kamailio-rtpengine/rootfs/etc/services.d/rtpengine/run
#!/usr/bin/with-contenv bash
set -euo pipefail

read_env() {
  [ -f "/run/s6/container_environment/$1" ] && cat "/run/s6/container_environment/$1" || true
}

IP="$(read_env ADVERTISED_IP)"
[ -n "${IP:-}" ] && [ "$IP" != "0.0.0.0" ] || { echo "FATAL: ADVERTISED_IP non valido"; exit 1; }

# Avvio "pulito": niente nft/iptables, niente daemon extra,
# sposto eventuali control-ports legacy su numeri innocui
# e tengo SOLO il control NG (UDP) che Kamailio usa.
exec rtpengine \
  --foreground \
  --interface="${IP}" \
  --listen-ng=127.0.0.1:22222 \
  --listen-udp=127.0.0.1:29223 \
  --listen-tcp=127.0.0.1:29224 \
  --port-min=49000 --port-max=50000 \
  --log-level=4 \
  --timeout=60 \
  --dtls-passive \
  --table=-1 \
  --nftables-chain=
